# vomsinfo/crontab
*/1 * * * * /usr/local/bin/vomsinfo > /dev/null 2>&1
